# infrastructure/maintenance/systemd/server-maintenance.service
[Unit]
Description=Server Maintenance (daily cleanup)
After=network.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/autoteka/deploy.env
EnvironmentFile=-/etc/autoteka/telegram.env
ExecStart=/usr/local/bin/autoteka maintenance
