# infrastructure/runtime/systemd/watch-changes.service.template
[Unit]
Description=Autoteka watch-changes service (git polling)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-${AUTOTEKA_OPTIONS_FILE}
ExecStart=/usr/local/bin/autoteka watch-changes
