# deploy/systemd/vue-app-deploy.service
[Unit]
Description=Autoteka Auto Deploy (git polling)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/vue-app/deploy.env
EnvironmentFile=-/etc/vue-app/telegram.env
ExecStart=/usr/local/bin/autoteka deploy
