# deploy/systemd/vue-app-deploy.service
[Unit]
Description=Vue App Auto Deploy (git polling)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/vue-app/deploy.env
WorkingDirectory=/opt/vue-app
ExecStart=/opt/vue-app/deploy/deploy.sh
