# deploy/runtime/systemd/autoteka.service
[Unit]
Description=Autoteka Frontend (docker compose up -d)
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/autoteka/deploy.env
ExecStart=/usr/local/bin/autoteka up
ExecStop=/usr/local/bin/autoteka down
RemainAfterExit=yes
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
